# shell/supervisor/broker.ini
[program:broker]
command=/usr/local/rocketmq-4.3.2/bin/mqbroker -n 0.0.0.0:9876
user=root
autostart=true
autorestart=true
stdout_logfile_maxbytes=20MB
stdout_logfile=/usr/local/rocketmq-4.3.2/logs/broker.log
stderr_logfile=/usr/local/rocketmq-4.3.2/logs/broker.log
