# shell/supervisor/nameserver.ini
[program:nameserver]
command=/usr/local/rocketmq-4.3.2/bin/mqnamesrv -n 0.0.0.0:9876
user=root
autostart=true
autorestart=true
stdout_logfile_maxbytes=20MB
stdout_logfile=/usr/local/rocketmq-4.3.2/logs/rocketmq.log
stderr_logfile=/usr/local/rocketmq-4.3.2/logs/rocketmq.log
