# shell/supervisor/redis.ini
[program:redis]
command=/usr/local/redis-3.2.5/src/redis-server /usr/local/redis-3.2.5/redis.conf
directory=/usr/local/redis-3.2.5
user=redis
autostart=true
autorestart=true
stdout_logfile_maxbytes=20MB
stdout_logfile=/usr/local/redis-3.2.5/6379.log
stderr_logfile=/usr/local/redis-3.2.5/6379.log