# shell/supervisor/tomcat.ini
[program:xdiamond]
command=/data/ftp/app/base/tomcat-diamond-new/bin/catalina.sh run
stdout_logfile=/data/ftp/app/base/tomcat-diamond-new/logs/catalina.out
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
