# shell/supervisor/zookeeper.ini
[program:zookeeper]
command=/usr/local/zookeeper/bin/zkServer.sh start-foreground /usr/local/zookeeper_data/2181/zoo.cfg
directory=/usr/local/zookeeper
user=zookeeper
autostart=true
autorestart=true
stopsignal=KILL
startsecs=10
startretries=999
log_stdout=true
log_stderr=true
logfile_maxbytes=20MB
logfile_backups=10d
stdout_logfile=/usr/local/zookeeper_data/2181/logs/zookeeper.out
stderr_logfile=/usr/local/zookeeper_data/2181/logs/zookeeper.out